# systemd/claudet-watchdog.timer
[Unit]
Description=claudet watchdog every 60s

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
