# systemd/claudet-watchdog.service
[Unit]
Description=claudet watchdog (one pass)

[Service]
Type=oneshot
ExecStart=/usr/bin/env node %h/.agent-os/bin/watchdog.mjs
