# launchd/com.claudet.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- claudet watchdog: one deterministic pass every 60s. Edit BOTH paths
     below (`which node` for the first; your home directory in the rest),
     then: launchctl load ~/Library/LaunchAgents/com.claudet.watchdog.plist
     Pause without unloading: touch ~/.agent-os/paused -->
<dict>
  <key>Label</key><string>com.claudet.watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/YOUR_USER/.agent-os/bin/watchdog.mjs</string>
  </array>
  <key>StartInterval</key><integer>60</integer>
  <key>StandardOutPath</key><string>/Users/YOUR_USER/.agent-os/launchd.log</string>
  <key>StandardErrorPath</key><string>/Users/YOUR_USER/.agent-os/launchd.log</string>
</dict>
</plist>
